# deb/src/main/resources/control/postinst
#!/bin/sh

# Ensure our service is enabled
mv /etc/bfm/bfmwatcher/bfm-app*.jar /etc/bfm/bfmwatcher/bfm-app.jar
systemctl enable bfm.service >/dev/null 2>&1

